# traffic_portal/test/supervisord/xvfb.conf
[program:xvfb]

command = Xvfb :1 -screen 0 1280x1024x24+32

stdout_logfile = /var/log/supervisor/%(program_name)s-out.log
stderr_logfile = /var/log/supervisor/%(program_name)s-err.log
